# rabbitmq/supervisor_rabbitmq.conf
[program:rabbitmq]
directory = /opt/mq/rabbitmq
command = /opt/mq/rabbitmq/sbin/rabbitmq-server start
autostart = true
autorestart = true
startsecs = 3
startretries = 3

