# rabbitmq/supervisor_rabbitmq.conf
[program:rabbitmq]
directory = /opt/mq/rabbitmq
command = /opt/mq/rabbitmq/sbin/rabbitmq-server -detached
autostart = true
autorestart = true
startsecs = 0
startretries = 3

