# glodontomcat/supervisor_sync_time.conf
[program:synctime]
directory = /opt/
command = /opt/sync_time.sh
autostart = true
autorestart = true
startsecs = 5
startretries = 3
# after all other operation
priority = 1
