# glodontomcat/supervisor_sync_resources.conf
[program:syncresources]
directory = /opt/
command = /opt/sync_resources.sh
autostart = true
autorestart = true
startsecs = 5
startretries = 3
# after all other operation
priority = 10000
