# wrapper_java/supervisor_wrapper.conf
[program:tomcat7]
directory = ${WRAPPER_HOME}
command = ${WRAPPER_HOME}/bin/startApp.sh console
autostart = true
autorestart = true
startsecs = 5
startretries = 3

