# wrapper_java/supervisor_wrapper.conf
[program:wrapper]
directory = /opt/wrapper
command = /opt/wrapper/bin/startApp.sh console
autostart = true
autorestart = true
startsecs = 5
startretries = 3

