# logstash/supervisor_logstash.conf
[program:logstash]
directory = /usr/share/logstash
command = logstash -f /logstash/conf/ -l /var/logs/logstash.log
autostart = true
autorestart = true
startsecs = 5
startretries = 3

