# codis-server/supervisor_codis.conf
[program:codis-server]
directory = /codis
command = /codis/bin/codis-server /codis/conf/redis.conf
autostart = true
autorestart = true
startsecs = 5
startretries = 3

